# deploy/systemd/ghostmon-agent.service
[Unit]
Description=GhostMon Agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ghostmon
Group=ghostmon
WorkingDirectory=/opt/ghostmon/agent
EnvironmentFile=-/etc/ghostmon/agent.env
ExecStart=/opt/ghostmon/agent/GhostMon.Agent
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077

[Install]
WantedBy=multi-user.target
